# init.d/teamviewerd
#!/sbin/runscript

name="teamviewerd"
description="Daemon for TeamViewer remote control program"
command="/opt/teamviewer8/tv_bin/${name}"
command_args="-d"
pidfile="/run/${name}.pid"

depend() {
    need net
}
